# systemd/hellfire-firewall.service
[Unit]
Description=Hellfire Firewall Service
Documentation=https://github.com/yourusername/hellfire
After=network.target hellfire-network.service
Requires=hellfire-network.service

[Service]
Type=oneshot
RemainAfterExit=yes

# Apply firewall rules from UCI config
ExecStart=/usr/local/bin/hf firewall apply --config-dir=/etc/config
ExecReload=/usr/local/bin/hf firewall reload --config-dir=/etc/config

# Flush rules on stop
ExecStop=/usr/local/bin/hf firewall flush --config-dir=/etc/config

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
